# home/.config/systemd/user/ssh-agent@.service
[Unit]
Description = SSH Agent Socket Forwarder

[Service]
ExecStart = pageant.exe
StandardInput = socket
StandardOutput = socket
StandardError = journal
