# home/.config/systemd/user/ssh-agent.socket
[Unit]
Description = SSH Agent Socket

[Socket]
ListenStream = /var/run/user/1000/ssh-agent.sock
Accept = Yes

[Install]
WantedBy = sockets.target
